# debian/fritzident.socket
[Unit]
Description=Fritz!Box parental control socket
Before=fritzident.service

[Socket]
ListenStream=14013

[Install]
WantedBy=sockets.target
