# debian/fritzident.socket
[Unit]
Description=Fritz!Box parental control socket

[Socket]
ListenStream=14013
Accept=true
