# debian/postinst
#!/bin/sh
# postinst script for fritzident
#
# see: dh_installdeb(1)

set -e
#set -x

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
    if [ -x /usr/sbin/update-inetd ]; then
	    update-inetd --remove 14013
	    update-inetd --add \
		    '14013\t\tstream\ttcp\tnowait\tnobody\t/usr/sbin/fritzident\t/usr/sbin/fritzident -l /var/log/fritzident.log'
            touch /var/log/fritzident.log
	    chown nobody:adm /var/log/fritzident.log
	    chmod 660 /var/log/fritzident.log
    fi
#    if [ -x /bin/systemctl ]; then
#	/bin/systemctl daemon-reload
 #       /bin/systemctl enable fritzident.socket
  #      /bin/systemctl start fritzident.socket
   #     /bin/systemctl enable fritzident.service
    #    /bin/systemctl start fritzident.service
    #fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
